# examples/ttn-mapper/ttn-mapper.service
  
[Unit]
Description=TTN Mapper service
After=gpsd.service

[Service]
ExecStart=/opt/lora/bin/ttn-mapper

[Install]
WantedBy=multi-user.target
